# etc/upstart/smartcaps.conf
# Config file for upstart

description "Smartcaps Daemon"

start on runlevel [2345]
stop on runlevel [!2345]

exec /usr/local/bin/smartcapsd
